# core/trait.sh.in
#!/bin/bash
# 
# chkconfig: 2345 20 80 
# description: Trait-o-matic service

USER=%USER%
SERVER=/home/$USER/trait-o-matic/core/server.py
SERVER_PID=/home/$USER/tmp/server.pid

start() {
	if [ -f "$SERVER_PID" ] && [ -n "$(ps -p $(< $SERVER_PID) -o cmd=)" ]; then
		echo "Trait-o-matic server already running"
	else
		echo "Starting Trait-o-matic server..."
		su $USER -c "python $SERVER &> /home/$USER/log/server.log" &
		echo $! > $SERVER_PID
	fi
}

stop() {
	echo "Stopping Trait-o-matic server..."
	kill $(< $SERVER_PID)
	rm $SERVER_PID
}

case "$1" in
	start)
		start
		;;
	stop)
		stop
		;;
	restart)
		stop
		start
		;;
	*)
		echo "Usage: $0 start|stop|restart"
		exit 1
		;;
esac
exit $?
